# car-stereo.service
# Car Stereo System Flask Backend Service
# Install to: /etc/systemd/system/car-stereo.service
#
# Usage:
#   sudo cp car-stereo.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable car-stereo
#   sudo systemctl start car-stereo
#

[Unit]
Description=Car Stereo System - Flask Backend
Documentation=https://github.com/your-repo/car_stereo_system
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/car_stereo_system
Environment=FLASK_ENV=production
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/car_stereo_system/venv/bin/python3 /home/pi/car_stereo_system/app.py
Restart=always
RestartSec=3
TimeoutStopSec=10

# Allow access to hardware
SupplementaryGroups=gpio i2c spi video audio plugdev

[Install]
WantedBy=multi-user.target

